# hardware/sparrow-bridge/skyguard-sparrow-bridge.service
[Unit]
Description=SkyGuard OS – Sparrow DroneID Bridge
Documentation=https://github.com/ghostop14/sparrow-wifi
After=network-online.target sparrow.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/skyguard/hardware/sparrow-bridge
EnvironmentFile=/home/pi/skyguard/hardware/sparrow-bridge/.env
ExecStart=/usr/bin/python3 /home/pi/skyguard/hardware/sparrow-bridge/bridge.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
