# hardware/droneID-scanner/skyguard-droneID-scanner.service
[Unit]
Description=SkyGuard OS — DroneID BLE Scanner (nRF52840)
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/skyguard/hardware/droneID-scanner
EnvironmentFile=/home/pi/skyguard/hardware/droneID-scanner/.env
ExecStart=/usr/bin/python3 /home/pi/skyguard/hardware/droneID-scanner/scanner.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
