# hardware/hackrf-bridge/skyguard-hackrf-bridge.service
[Unit]
Description=SkyGuard OS HackRF Sweep Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/admin/skyguard/hardware/hackrf-bridge
EnvironmentFile=/home/admin/skyguard/hardware/hackrf-bridge/.env
ExecStart=/home/admin/hackrf_sweep/venv/bin/python3 /home/admin/skyguard/hardware/hackrf-bridge/bridge.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
